# etc/profile.d/pythoncerts.sh
version https://git-lfs.github.com/spec/v1
oid sha256:624686abf4b3b15d405dcfb612f5b6b969f2953213efe2a656d71346843b7527
size 137
